# scripts/systemd/k3s-start-day.service
[Unit]
Description=K3s Day Mode - Start K3s
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s-start.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
